# com.jancok.smbcheck.python.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jancok.smbcheck.python</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jancok/Documents/bash_script/smb_check.py</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jancok/smb_monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jancok/smb_monitor.log</string>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>